# vscpd.startup_script_fedora
#!/bin/bash
#
# vscpd       This shell script takes care of starting and stopping
#             standalone vscpd.
#
# chkconfig: - 60 50
# description: Vscpd is an automation daemon, which is the program \
#              that handles controldrivers and clients.
# processname: vsftpd
# config: /etc/vscp/vscp.conf or /etc/vscp.conf

# Source function library.
. /etc/rc.d/init.d/functions

[ -x /usr/local/bin/vscpd ] || exit 0

RETVAL=0
prog="vscpd"

start() {
        # Start daemons.

        if [ -d /etc/vscp ] ; then
                for i in `ls /etc/vscp/*.conf`; do
                        site=`basename $i .conf`
                        echo -n $"Starting $prog for $site: "
                        /usr/local/bin/vscpd $i &
                        RETVAL=$?
                        [ $RETVAL -eq 0 ] && {
                           touch /var/lock/subsys/$prog
                           success $"$prog $site" 
                        } 
                        echo
                done
        else
                RETVAL=1
        fi
        return $RETVAL
}

stop() {
        # Stop daemons.
        echo -n $"Shutting down $prog: "
        killproc $prog
        RETVAL=$?
        echo
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/$prog
        return $RETVAL
}

# See how we were called.
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|reload)
        stop
        start
        RETVAL=$?
        ;;
  condrestart)
        if [ -f /var/lock/subsys/$prog ]; then
            stop
            start
            RETVAL=$?
        fi
        ;;
  status)
        status $prog
        RETVAL=$?
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|condrestart|status}"
        exit 1
esac

exit $RETVAL
