# service/svscan/run
#!/bin/sh

/usr/sbin/vscpd -c /etc/vscp/vscpd.json
