# service/upstart/vscpd.conf
description "vscpd daemon"
author "Ake Hedman <akhe@vscp.org>"

start on net-device-up

respawn

exec /usr/local/sbin/vscpd -c /etc/vscp/vscpd.conf
