# config_examples/vscpd.startup_script_carambola
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
START=10
STOP=15

start() {        
    echo start
    /usr/bin/vscpd
}                 
                   
stop() {          
    echo stop
    killall vscpd 
}

restart() {
    echo restart
    killall vscpd
    /usr/bin/vscpd
}
